# network-check/wifi_monitor/wifi-monitor.service
[Unit]
Description=WiFi Network Monitor Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/wifi_monitor
Environment=PYTHONPATH=/home/pi/wifi_monitor
ExecStart=/home/pi/wifi_monitor/wifi_monitor_env/bin/python /home/pi/wifi_monitor/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 로그 파일 위치
StandardOutput=append:/home/pi/wifi_monitor/logs/wifi_monitor.log
StandardError=append:/home/pi/wifi_monitor/logs/wifi_monitor_error.log

[Install]
WantedBy=multi-user.target
